# platform/broadcom/sonic-platform-modules-cel/debian/platform-modules-seastone2.postinst
depmod -a
systemctl enable platform-modules-seastone2.service
systemctl start platform-modules-seastone2.service