# platform/broadcom/sonic-platform-modules-cel/debian/platform-modules-dx010.postinst
depmod -a
systemctl enable platform-modules-dx010.service
systemctl enable fancontrol.service

systemctl start platform-modules-dx010.service
systemctl start fancontrol.service

/usr/local/bin/platform_api_mgnt.sh install

