# platform/broadcom/sonic-platform-modules-cel/debian/platform-modules-haliburton.postinst
depmod -a
systemctl enable platform-modules-haliburton.service
systemctl start platform-modules-haliburton.service

systemctl enable cpu_wdt.service
systemctl start cpu_wdt.service
