# platform/broadcom/sonic-platform-modules-cel/debian/platform-modules-silverstone.postinst
depmod -a
systemctl enable platform-modules-silverstone.service
systemctl start platform-modules-silverstone.service

/usr/local/bin/platform_api_mgnt.sh install
