# platform/broadcom/sonic-platform-modules-cel/debian/platform-modules-belgite.postinst
depmod -a
/usr/local/bin/platform_api_mgnt.sh install
systemctl enable pddf-platform-init.service
systemctl start pddf-platform-init.service
systemctl enable belgite-pddf-platform-monitor.service
systemctl start belgite-pddf-platform-monitor.service
