# platform/nephos/sonic-platform-modules-accton/debian/sonic-platform-accton-as7116-54x.postinst
depmod -a
systemctl enable as7116-platform-init.service
systemctl start  as7116-platform-init.service

/usr/local/bin/platform_api_mgnt.sh install
