# platform/broadcom/sonic-platform-modules-ingrasys/s8810-32q/service/qsfp-monitor.service
[Unit]
Description= This QSFP Monitor service is to setup QSFP SI.
Requires=s8810-32q-monitor.service
After=s8810-32q-monitor.service

[Service]
ExecStart=/usr/sbin/qsfp_monitor.sh
KillSignal=SIGKILL
SuccessExitStatus=SIGKILL

# Resource Limitations
LimitCORE=infinity

[Install]
WantedBy=multi-user.target
