# platform/broadcom/sonic-platform-modules-juniper/debian/sonic-platform-juniper-qfx5210.postinst
systemctl enable qfx5210-platform-init.service
systemctl start qfx5210-platform-init.service
