# platform/broadcom/sonic-platform-modules-dell/s5232f/systemd/platform-modules-s5232f.service
[Unit]
Description=Dell S5232f Platform modules
Before=pmon.service
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/usr/local/bin/s5232f_platform.sh init
ExecStop=/usr/local/bin/s5232f_platform.sh deinit
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
