# files/image_config/rsyslog/rsyslog-config.service
[Unit]
Description=Update rsyslog configuration

[Service]
Type=oneshot
ExecStart=/usr/bin/rsyslog-config.sh

[Install]
WantedBy=multi-user.target
