# platform/barefoot/sonic-platform-modules-bfn-newport/debian/postinst
#!/bin/sh
set -e
depmod -a 
systemctl enable bfn-newport.service
systemctl start bfn-newport.service
#DEBHELPER#
