# platform/marvell-armhf/sonic-platform-nokia/debian/sonic-platform-nokia-7215.postinst
#!/bin/sh
# postinst script for sonic-platform-nokia-7215
#
# see: dh_installdeb(1)

sh /usr/sbin/nokia-7215_plt_setup.sh
systemctl enable nokia-7215init.service
systemctl start nokia-7215init.service

systemctl enable cpu_wdt.service
systemctl start cpu_wdt.service

exit 0

