# platform/centec-arm64/tsingma-bsp/debian/tsingma-bsp.postinst
systemctl enable tsingma-bsp.service
systemctl start tsingma-bsp.service
