# platform/broadcom/sonic-platform-modules-accton/debian/sonic-platform-accton-as7816-64x.prerm
# Disable pddf-platform-init, monitor, monitor-fan, monitor-psu

systemctl stop pddf-platform-init.service
systemctl disable pddf-platform-init.service
systemctl stop as7816-pddf-platform-monitor.service
systemctl disable as7816-pddf-platform-monitor.service

