# files/image_config/fstrim/fstrim.service
[Unit]
Description=Discard unused blocks

[Service]
Type=oneshot
ExecStartPre=/usr/local/bin/log_ssd_health
ExecStart=/sbin/fstrim -av
