# files/image_config/process-reboot-cause/process-reboot-cause.service
[Unit]
Description=Reboot cause determination service
After=rc-local.service

[Service]
Type=simple
ExecStart=/usr/bin/process-reboot-cause
