# platform/s3ip-sysfs/debian/s3ip-sysfs.prerm
systemctl stop s3ip-sysfs
systemctl disable s3ip-sysfs
