# platform/s3ip-sysfs/debian/s3ip-sysfs.postinst
systemctl daemon-reload
systemctl start s3ip-sysfs
systemctl enable s3ip-sysfs
