# platform/barefoot/sonic-platform-modules-ingrasys/debian/sonic-platform-ingrasys-s9180-32x.postrm
# Automatically added by dh_systemd_start
if [ -d /run/systemd/system ]; then
    systemctl --system daemon-reload >/dev/null || true
    fi
# End automatically added section
# Automatically added by dh_installinit
if [ "$1" = "purge" ] ; then
    update-rc.d s9180-32x-monitor remove >/dev/null
    update-rc.d bfn remove >/dev/null
    update-rc.d qsfp-monitor remove >/dev/null
fi


# In case this system is running systemd, we make systemd reload the unit files
# to pick up changes.
if [ -d /run/systemd/system ] ; then
    systemctl --system daemon-reload >/dev/null || true
fi
# End automatically added section
# Automatically added by dh_systemd_enable
if [ "$1" = "remove" ]; then
    if [ -x "/usr/bin/deb-systemd-helper" ]; then
        deb-systemd-helper mask s9180-32x-monitor.service >/dev/null
        deb-systemd-helper mask bfn.service >/dev/null
        deb-systemd-helper mask qsfp-monitor.service >/dev/null
    fi
fi

if [ "$1" = "purge" ]; then
    if [ -x "/usr/bin/deb-systemd-helper" ]; then
        deb-systemd-helper purge s9180-32x-monitor.service >/dev/null
        deb-systemd-helper unmask s9180-32x-monitor.service >/dev/null
        deb-systemd-helper purge bfn.service >/dev/null
        deb-systemd-helper unmask bfn.service >/dev/null
        deb-systemd-helper purge qsfp-monitor.service >/dev/null
        deb-systemd-helper unmask qsfp-monitor.service >/dev/null
    fi
fi
# Generate kernel modules.dep and map files for remove eeprom_mb.
depmod -a || true
# End automatically added section

