# platform/broadcom/sonic-platform-modules-dell/s6100/systemd/s6100-ssd-monitor.timer
[Unit]
Description=Dell S6100 SSD monitoring poller timer
DefaultDependencies=no
After=pmon.service

[Timer]
OnBootSec=5min
OnUnitActiveSec=60min

[Install]
WantedBy=timers.target

