# platform/centec/sonic-platform-modules-embedway/debian/platform-modules-embedway-es6220.postinst
depmod -a
systemctl enable es6220_platform.service
systemctl start es6220_platform.service
