# files/image_config/corefile_uploader/core_uploader.service
[Unit]
Description=Host core file uploader daemon
Requires=updategraph.service
After=updategraph.service

[Service]
Type=simple
ExecStart=/usr/bin/core_uploader.py

[Install]
WantedBy=multi-user.target
