# platform/centec/sonic-platform-modules-v682/debian/platform-modules-v682-48y8c-d.postinst
systemctl enable 48y8c_d_platform.service
systemctl start 48y8c_d_platform.service
