# platform/centec/sonic-platform-modules-v682/debian/platform-modules-v682-48y8c.postinst
systemctl enable 48y8c_platform.service
systemctl start 48y8c_platform.service
