# platform/centec/sonic-platform-modules-v682/debian/platform-modules-v682-48x8c.postinst
systemctl enable 48x8c_platform.service
systemctl start 48x8c_platform.service
