# platform/centec-arm64/sonic-platform-modules-e530/debian/platform-modules-e530-24x2q.postinst
systemctl enable 24x2q_platform.service
systemctl start 24x2q_platform.service
