# platform/centec-arm64/sonic-platform-modules-e530/debian/platform-modules-e530-48t4x-p.postinst
systemctl enable 48t4x_p_platform.service
systemctl start 48t4x_p_platform.service
