# platform/centec-arm64/sonic-platform-modules-e530/debian/platform-modules-e530-48s4x.postinst
systemctl enable 48s4x_platform.service
systemctl start 48s4x_platform.service
