# platform/marvell-armhf/sonic-platform-nokia/7215/service/cpu_wdt.service
[Unit]
Description=CPU WDT
After=nokia-7215init.service
[Service]
ExecStart=/usr/local/bin/cpu_wdt.py

[Install]
WantedBy=multi-user.target
