# platform/barefoot/sonic-platform-modules-ufispace/debian/sonic-platform-ufispace-s9180-32x.postinst
depmod -a
systemctl enable pddf-platform-init.service
systemctl start pddf-platform-init.service
systemctl enable platform-swss-restart.service
systemctl start --no-block platform-swss-restart.service
