# files/image_config/pcie-check/pcie-check.service
[Unit]
Description=Check the PCIe device presence and status
After=rc.local.service

[Service]
Type=simple
ExecStart=/usr/bin/pcie-check.sh
