# platform/broadcom/sonic-platform-modules-ruijie/debian/platform-modules-ruijie-b6510-48vs8cq.init
#!/bin/bash

### BEGIN INIT INFO
# Provides:          setup-board
# Required-Start:
# Required-Stop:
# Should-Start:
# Should-Stop:
# Default-Start:     S
# Default-Stop:      0 6
# Short-Description: Setup B6510-48VS8CQ board.
### END INIT INFO

case "$1" in
start)
    echo -n "Setting up board... "
	
    /usr/local/bin/b6510-48vs8cq_platform.sh init

    echo "done."
    ;;

stop)
    /usr/local/bin/b6510-48vs8cq_platform.sh deinit
    echo "done."

    ;;

force-reload|restart)
    echo "Not supported"
    ;;

*)
    echo "Usage: /etc/init.d/platform-modules-b6510-48vs8cq.init {start|stop}"
    exit 1
    ;;
esac

exit 0
