# platform/centec-arm64/sonic-platform-modules-ragile/ra-b6010-48gt4x/service/ra-b6010-48gt4x_platform.service
[Unit]
Description=Centec modules init
After=local-fs.target
Before=syncd.service

[Service]
Type=oneshot
ExecStart=-/etc/init.d/platform-modules-ra-b6010-48gt4x start
ExecStop=-/etc/init.d/platform-modules-ra-b6010-48gt4x stop
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
