# platform/broadcom/sonic-platform-modules-dell/s5212f/systemd/platform-modules-s5212f.service
[Unit]
Description=Dell S5212f Platform modules
Before=pmon.service
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/usr/local/bin/s5212f_platform.sh init
ExecStop=/usr/local/bin/s5212f_platform.sh deinit
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
