# src/sonic-host-services-data/debian/sonic-host-services-data.process-reboot-cause.service
[Unit]
Description=Retrieve the reboot cause from the history files and save them to StateDB
Requires=database.service
After=database.service

[Service]
Type=simple
ExecStart=/usr/local/bin/process-reboot-cause
