# platform/centec-arm64/sonic-platform-modules-ragile/common/service/fancontrol.service
[Unit]
Description= Fancontrol for Device.
After=privatenetwork.service
DefaultDependencies=no

[Service]
ExecStart=/usr/local/bin/fancontrol.py start
Restart=on-failure

[Install]
WantedBy=multi-user.target

