# platform/innovium/sonic-platform-modules-cel/debian/platform-modules-silverstone-x.postinst
for i in $(seq 1 5); do
    ret=$(depmod -a)
    if [ -n "$ret" ]; then
        break
    fi
    sleep 2
    echo "run depmod again"
done
systemctl enable platform-modules-silverstone-x.service
systemctl start platform-modules-silverstone-x.service
systemctl enable fancontrol.service
systemctl start fancontrol.service

/usr/local/bin/platform_api_mgnt.sh install
/etc/init.d/fancontrol.sh install
