# platform/innovium/sonic-platform-modules-cel/debian/platform-modules-midstone-200i.postinst
depmod -a
systemctl enable platform-modules-midstone-200i.service
systemctl start platform-modules-midstone-200i.service