# platform/innovium/sonic-platform-modules-cel/debian/platform-modules-midstone-100x.postinst
depmod -a
systemctl enable platform-modules-midstone-100x.service
systemctl start platform-modules-midstone-100x.service

/usr/local/bin/platform_api_mgnt.sh install
