# platform/broadcom/sonic-platform-modules-dell/debian/platform-modules-z9100.postinst
# postinst script for Z9100

# Enable Dell-Z9100-platform-service
depmod -a 
systemctl enable platform-modules-z9100.service
systemctl start platform-modules-z9100.service


systemctl enable z9100-lpc-monitor.service
systemctl start z9100-lpc-monitor.service

#DEBHELPER#
