# platform/broadcom/sonic-platform-modules-dell/debian/platform-modules-s6000.postinst
# postinst script for S6000

# Enable Dell-S6000-platform-service
depmod -a 
systemctl enable platform-modules-s6000.service
systemctl start platform-modules-s6000.service

systemctl enable fancontrol.service
systemctl start fancontrol.service


#DEBHELPER#
