# platform/broadcom/sonic-platform-modules-dell/debian/platform-modules-s5212f.postinst
# postinst script for S5212f

# Enable Dell-S5212f-platform-service
depmod -a
systemctl enable platform-modules-s5212f.service
systemctl start platform-modules-s5212f.service


#DEBHELPER#

