# platform/broadcom/sonic-platform-modules-dell/debian/platform-modules-n3248te.postinst
# postinst script for N3248TE

# Enable Dell-N3248TE-platform-service
depmod -a
systemctl enable platform-modules-n3248te.service
systemctl start platform-modules-n3248te.service
#DEBHELPER#
