# platform/broadcom/sonic-platform-modules-dell/debian/platform-modules-n3248pxe.postinst
# postinst script for n3248pxe

# Enable Dell-n3248pxe-platform-service
depmod -a
systemctl enable platform-modules-n3248pxe.service
systemctl start platform-modules-n3248pxe.service
#DEBHELPER#
