# platform/broadcom/sonic-platform-modules-dell/debian/platform-modules-z9432f.postinst
# postinst script for Z9432f

# Enable Dell-Z9432f-platform-service
depmod -a
systemctl enable platform-modules-z9432f.service
systemctl start platform-modules-z9432f.service


#DEBHELPER#

