# platform/broadcom/sonic-platform-modules-dell/debian/platform-modules-z9332f.postinst
# postinst script for Z9332f

# Enable Dell-Z9332f-platform-service
depmod -a
systemctl enable platform-modules-z9332f.service
systemctl start platform-modules-z9332f.service


#DEBHELPER#

