# platform/broadcom/sonic-platform-modules-dell/debian/platform-modules-z9100.postinst
# postinst script for Z9100

# Enable fstrim
systemctl enable fstrim.timer
systemctl start fstrim.timer

#DEBHELPER#
