# platform/broadcom/sonic-platform-modules-dell/debian/platform-modules-s5224f.postinst
# postinst script for S5224f

# Enable Dell-S5224f-platform-service
depmod -a
systemctl enable platform-modules-s5224f.service
systemctl start platform-modules-s5224f.service


#DEBHELPER#

