# dockers/docker-fpm-gobgp/supervisord.conf
[supervisord]
nodaemon=true

[program:gobgpd]
command=/usr/sbin/gobgpd -p -f /etc/gobgp/gobgpd.conf -r
priority=1
