# files/image_config/interfaces/interfaces-config.service
[Unit]
Description=Update interfaces configuration
Before=database.service

[Service]
Type=oneshot
ExecStart=/usr/bin/interfaces-config.sh

[Install]
WantedBy=multi-user.target

