# platform/innovium/sonic-platform-modules-cameo/debian/sonic-platform-cameo-esc601-32q.postinst
#!/bin/sh

depmod -a
systemctl enable esc601-platform-init.service
systemctl start esc601-platform-init.service


/usr/local/bin/cameo_esc601_platform.sh
