# platform/innovium/sonic-platform-modules-cameo/debian/sonic-platform-cameo-esqc610-56sq.postinst
#!/bin/sh

depmod -a

systemctl enable esqc610-platform-init.service
systemctl start esqc610-platform-init.service
/usr/local/bin/cameo_esqc610_platform.sh