# platform/innovium/sonic-platform-modules-cameo/debian/sonic-platform-cameo-esc602-32q.postinst
#!/bin/sh

depmod -a
systemctl enable esc602-platform-init.service
systemctl start esc602-platform-init.service


/usr/local/bin/cameo_esc602_platform.sh