# platform/innovium/sonic-platform-modules-cameo/debian/sonic-platform-cameo-escc601-32q.postinst
#!/bin/sh

depmod -a
systemctl enable escc601-platform-init.service
systemctl start escc601-platform-init.service


/usr/local/bin/cameo_escc601_platform.sh