# platform/centec-arm64/sonic-platform-modules-ragile/debian/platform-modules-ra-b6010-48gt4x.postinst
systemctl enable ra-b6010-48gt4x_platform.service
systemctl start ra-b6010-48gt4x_platform.service
systemctl enable device_i2c.service
systemctl start device_i2c.service
systemctl enable privatenetwork.service
systemctl start privatenetwork.service
systemctl enable fancontrol.service
systemctl start fancontrol.service
