# platform/centec/sonic-platform-modules-e582/debian/platform-modules-e582-48x2q4z.postinst
depmod -a
systemctl enable 48x2q4z_platform.service
systemctl start 48x2q4z_platform.service
