# platform/centec/sonic-platform-modules-e582/debian/platform-modules-e582-48x6q.postinst
depmod -a
systemctl enable 48x6q_platform.service
systemctl start 48x6q_platform.service
