# platform/innovium/sonic-platform-modules-wistron/sw-to3200k/service/to3200k-platform.service
[Unit]
Description=Wistron SW-TO3200K Platform Monitoring service
Before=pmon.service
After=sysinit.target
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/usr/local/bin/platform_sync
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
