# platform/broadcom/sonic-platform-modules-cel/haliburton/systemd/cpu_wdt.service
[Unit]
Description=CPU WDT
After=platform-modules-haliburton.service
Requires=platform-modules-haliburton.service

[Service]
ExecStart=-/usr/local/bin/cpu_wdt start

[Install]
WantedBy=multi-user.target
