# files/image_config/cron.d/core_cleanup
# Attempts to clean up core files every 2 hours
0 */2 * * * root /usr/bin/core_cleanup.py > /dev/null 2>&1

