# files/image_config/cron.d/logrotate
# Attempt to rotate system logs once per minute
* * * * * root /usr/sbin/logrotate /etc/logrotate.d/rsyslog
