# files/image_config/updategraph/updategraph.service
[Unit]
Description=download minigraph from graph service
Before=database.service

[Service]
Type=oneshot
ExecStart=/usr/bin/updategraph

[Install]
WantedBy=multi-user.target
