# platform/broadcom/sonic-platform-modules-ufispace/debian/sonic-platform-ufispace-s8901-54xc.prerm
systemctl stop pddf-platform-init.service
systemctl disable pddf-platform-init.service