# platform/broadcom/sonic-platform-modules-ufispace/debian/sonic-platform-ufispace-s7801-54xs.postinst
depmod -a
systemctl enable pddf-platform-init.service
systemctl start pddf-platform-init.service
