# files/build_templates/enable_counters.service
[Unit]
Description=Enable SONiC counters

[Service]
Type=oneshot
ExecStart=/usr/bin/counterpoll queue enable
RemainAfterExit=yes
