# files/build_templates/sonic-delayed.target
[Unit]
Description=SONiC delayed services target.

[Install]
WantedBy=multi-user.target
