# dockers/docker-database/supervisord.conf
[supervisord]
logfile_maxbytes=1MB
logfile_backups=2
nodaemon=true

[program:redis-server]
command=/usr/bin/redis-server /etc/redis/redis.conf
priority=1
autostart=true
autorestart=false
stdout_logfile=syslog
stderr_logfile=syslog
